# .github/systemd/qwerty-midi.service
[Unit]
Description=QWERTY MIDI Keyboard Monitor
Documentation=https://github.com/yourusername/qwerty-midi-keyboard
After=network.target sound.target

[Service]
Type=simple
ExecStart=REPLACE_WITH_PYTHON_PATH REPLACE_WITH_SCRIPT_PATH
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=default.target

